# contrib/debian/postinst
#!/bin/sh
set -e

# Reload systemd units and enable/start exporter
if command -v systemctl >/dev/null 2>&1; then
    systemctl daemon-reload || true
    systemctl enable policy-rate-limit.service || true
    systemctl restart policy-rate-limit.service || true
fi
exit 0
